# roles/next.js/templates/nextjs.service.j2
[Unit]
Description=Next.js Frontend
After=network.target

[Service]
Type=simple
User={{ app_user }}
WorkingDirectory={{ nextjs_path }}

ExecStart=/usr/bin/pnpm start

Restart=always
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
